# stack/units/early-services.target
[Unit]
Requires=docker.service
After=docker.service

[Install]
WantedBy=services.target
